# cerbo/service/log/run
#!/bin/sh
mkdir -p /var/log/cerbo-p1-bridge
exec multilog t s25000 n4 /var/log/cerbo-p1-bridge
